# bare-metal/ansible/roles/loki/templates/loki.service.j2
[Unit]
Description=Loki
After=network.target

[Service]
Type=simple
ExecStart={{ loki_bin }} -config.file={{ loki_config }}
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
