# bare-metal/ansible/roles/content-service/templates/content-service.service.j2
[Unit]
Description=EdTech Content Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User={{ app_user }}
Group={{ app_group }}
WorkingDirectory={{ service_dir }}/repo/content-service
EnvironmentFile={{ service_dir }}/.env
ExecStart=/usr/bin/node {{ service_dir }}/repo/content-service/dist/index.js
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=content-service
TimeoutStopSec=30s
KillMode=mixed
KillSignal=SIGTERM
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths={{ service_dir }} {{ log_dir }}

[Install]
WantedBy=multi-user.target
